# transport/coap/target/control/deb/postinst
#!/bin/sh

set -e

mkdir -m 0755 -p /var/log/tb-coap-transport
chown -R thingsboard: /var/log/tb-coap-transport
chown -R thingsboard: /usr/share/tb-coap-transport
systemctl --no-reload enable tb-coap-transport.service >/dev/null 2>&1 || :

exit 0
